# deploy/setup/InTheZone.service
[Unit]
Description = InTheZone game server

[Service]
ExecStart=/home/ubuntu/devrepo/deploy/startServer.sh

[Install]
WantedBy=multi-user.target

